# .config/systemd/user/mailnotif.timer
[Unit]
Description=Mail notification timer

[Timer]
OnUnitActiveSec=3m
Unit=mbsync.service

[Install]
WantedBy=timers.target