# .config/systemd/user/mailnotif.timer
[Unit]
Description=Mail notification timer

[Timer]
OnUnitActiveSec=1m
Unit=mailnotif.service

[Install]
WantedBy=timers.target