# .config/systemd/user/mbsync.timer
[Unit]
Description=Mailbox synchronization timer

[Timer]
OnBootSec=2m
OnUnitActiveSec=2h
Unit=mbsync.service

[Install]
WantedBy=timers.target