# .config/systemd/user/polybar-mail.service
[Unit]
Description=Polybar mail count

[Service]
Type=idle
ExecStart=%h/.local/bin/polybar/mail-count

[Install]
WantedBy=default.target