# .config/systemd/user/syncthing-perm.timer
[Unit]
Description=Syncthing permission setup timer

[Timer]
OnBootSec=2m
OnUnitActiveSec=1h
Unit=syncthing-perm

[Install]
WantedBy=timers.target