# source/syncthing/vts-syncthing-init.service
# This systemd unit file is automatically generated by Org-babel and do not edit this file manually.
[Unit]
Description=Syncthing permission init
After=network.target home-syncthing.mount

[Service]
Type=simple
ExecStart=/usr/local/bin/vts-syncthing-init
CPUWeight=20
CPUQuota=60%
IOWeight=20
# MemorySwapMax=0

[Install]
WantedBy=multi-user.target
